# com.enconvo.telegram-adapter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.enconvo.telegram-adapter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zanearcher/Downloads/vibe-coding/enconvo_cli/scripts/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zanearcher/.local/share/enconvo-telegram-adapter</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>NetworkState</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zanearcher/Library/Logs/enconvo-telegram-adapter.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zanearcher/Library/Logs/enconvo-telegram-adapter-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
